# PKGBUILDs/ovos-enclosure-rpi4-mark2/rpi4-mark2.install
# This script is only meant to be ran on the actual RPi4 device.
post_install() {
    systemctl enable systemd-timesyncd 
    systemctl enable pulseaudio.service
    systemctl disable getty@tty1
}

post_upgrade() {
    post_install
}